# debian/retaliq-domain.service
[Unit]
Description=Retaliq domain helper
After=network.target

[Service]
Type=simple
# explicit subcommand keeps unit safe from CLI default changes
ExecStart=/usr/bin/retaliq-domain serve -config /etc/retaliq-domain/config.conf
Restart=on-failure

[Install]
WantedBy=multi-user.target
